# warmshower-rate-mirrors-pkg/cachyos-rate-mirrors.install
post_install() {
    systemctl daemon-reload
    systemctl enable warmshower-rate-mirrors.timer
}

post_upgrade() {
    systemctl daemon-reload
    systemctl enable warmshower-rate-mirrors.timer
}

pre_remove() {
    systemctl disable warmshower-rate-mirrors.timer
}
